# modules/ansible-roles/home_assistant/templates/hass.service.j2
[Unit]
Description=Home Assistant
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hass

[Install]
WantedBy=multi-user.target